# packaging/systemd/opp_ci-serve.service
[Unit]
Description=opp_ci coordinator (web UI + API + scheduler)
Documentation=https://github.com/omnetpp/opp_ci
After=network-online.target postgresql.service
Wants=network-online.target
PartOf=opp_ci.target

[Service]
Type=simple
User=opp_ci
Group=opp_ci
# Read access to the system journal so the web UI's Logs pages can tail
# the serve and worker units (`journalctl -u …`). Scoped to this process —
# no change to the opp_ci account's global groups.
SupplementaryGroups=systemd-journal
WorkingDirectory=/var/lib/opp_ci
EnvironmentFile=/etc/opp_ci/opp_ci.env
EnvironmentFile=-/etc/opp_ci/serve.env
# Put the venv's bin/ on PATH so subprocess calls (opp_env install, opp_env
# run, ...) resolve without per-call-site path tricks. Mirrors what
# `source /opt/opp_ci/setenv` does for an interactive shell.
Environment="PATH=/opt/opp_ci/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/opp_ci/.venv/bin/opp_ci serve
Restart=on-failure
RestartSec=5s

# Hardening — safe for serve (no Nix, no podman).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/opp_ci

[Install]
WantedBy=opp_ci.target
